# addons/ff-vdr/source/smp-affinity.service.d/override.conf
[Unit]
Before=vdr.service vdr-frontend.service

[Install]
WantedBy=vdr.service vdr-frontend.service
