# addons/ff-vdr/source/locale.service.d/override.conf
[Unit]
Before=vdr.service vdr-frontend.service
