# addons/ff-vdr/source/kodi-halt.service.d/override.conf
[Unit]
After=vdr.service vdr-frontend.service
